# summarizer-bot.service.example
[Unit]
Description=Discord Summarizer Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=<YOUR_USER>
WorkingDirectory=/opt/summarizer-bot
ExecStart=/opt/summarizer-bot/summarizer-bot
Restart=on-failure
RestartSec=5
EnvironmentFile=/opt/summarizer-bot/.env

[Install]
WantedBy=multi-user.target
